# update-motd.d/05-legal
#!/bin/sh
# show the legal page upon login
[ -f /etc/legal ] && /bin/cat /etc/legal
