# update-motd.d/40-hostname
#!/bin/sh
/bin/hostname smtp.clouderg.com
/bin/echo smtp.clouderg.com > /etc/hostname
